# .jekyll-cache/Jekyll/Cache/Jekyll--Converters--Markdown/06/4ad60a6d679082f3862782bc828a6e4d763268ed34a57cc86dd23a4f69b610
I"¦<p>ngrinderì´ ë³´ë‹ˆ.. ë°±ê·¸ë¼ìš´ë“œ í˜•ì‹ìœ¼ë¡œ ë°–ì— ëŒë¦¬ì§€ë¥¼ ëª»í•˜ëŠ” í˜•ì‹ìœ¼ë¡œ ì œê³µì´ ëœë‹¤..</p>

<p>&amp;ë¶™ì—¬ì„œ ëŒë¦¬ë©´ ë˜ê¸´ í•˜ì§€ë§Œ ì‚´ì§ ì•„ì‰¬ìš´ ê²½ìš°ë¼.. screenì„ ì´ìš©í•´ì„œ ì•„ë˜ì™€ ê°™ì´ ë°ëª¬ í–‰íƒœë¡œ ëŒë¦´ìˆ˜ ìˆë‹¤..</p>

<p>screenëŠ” ì„¤ì¹˜ê°€ ë˜ì–´ ìˆì–´ì•¼ í•œë‹¤..</p>

<pre class="brush:bash">#!/bin/bash
#
# ngrinder_agent         Start up the ngrinder_agent daemon
#
# processname: ngrinder_agent
### BEGIN INIT INFO
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
#                    This service starts up the ngrinder_agent server daemon.
### END INIT INFO


DAEMON="ngrinder_agent"
IDENT="ngrinder_agent"
PREFIX="/home/ngrinder-core-3.2.3"
RUNFILE="/var/run/ngrinder_agent.run"


start() {
        [ -e ${RUNFILE} ] &amp;&amp; stop
        touch $RUNFILE
        echo -n "starting $DAEMON"
        screen -d -m -S ${IDENT} ${PREFIX}/run_agent.sh &amp;
        echo -n "."
        echo ""
}

stop() {
        echo "stopping $DAEMON"
        screen -list | grep \.${IDENT} | cut -d\. -f1 | xargs -r kill -9
        screen -wipe &gt; /dev/null 2&gt;&amp;1
        rm -f ${RUNFILE} 2&gt;/dev/null
}

case "$1" in
        start)  start
                ;;
        stop)   stop
                ;;
        status)
          chk=`screen -list | grep \.${IDENT} | awk '{print $1}' | wc -l`
          ppid=`screen -list | grep ngrinder_agent | awk '{print $0}' | awk -F\. '{print $1}'| xargs`
                if [ $chk == 1 ]
                        then
                          echo "$DAEMON is (pid $ppid) running.."
                         exit 0
                        else
                          echo "$DAEMON is not running.."
                         exit 1
                fi
                ;;
        restart)stop
                start
                ;;
        *)      echo -n "usage:\n\t$0 (start|stop|status|restart)\n"
                exit 1
                ;;
esac
exit 0

</pre>
:ET